# crontab
* * * * * /bin/bash -l -c "cd /usr/src && ruby runner.rb"